# milkv-arduino/lib/systemd/system/milkv-arduino.service
[Unit]
Description=Milk-V Duo Arduino Upload Daemon

[Service]
ExecStart=/usr/libexec/mdaud
Restart=always

[Install]
WantedBy=multi-user.target